# .config/autostart/keynav.desktop
[Desktop Entry]
Exec=/usr/bin/keynav
Icon=
Name=keynav
Path=
Terminal=False
Type=Application
